# bin/merge.service
[Unit]
Description=Merge Sort in C++
After=network.target

[Service]
WorkingDirectory=/bin/
ExecStart=/usr/bin/main
User=Merge
Restart=always

[Install]
WantedBy=multi-user.target