# bin/merge.service
[Unit]
Description=Merge Sort in Python
After=network.target

[Service]
WorkingDirectory=/bin/
ExecStart=/usr/bin/python3 /bin/main.py
User=merge
Restart=always

[Install]
WantedBy=multi-user.target
