# services/DEBIAN/postinst
#!/bin/sh

set -e
echo "Starting postinstall"

adduser --system --home /var/lib/merge merge

systemctl daemon-reload
systemctl enable merge
systemctl start merge