# services/DEBIAN/postinst
#!/bin/sh

set -e
echo "Starting postinstall"

adduser --system merge

systemctl daemon-reload
systemctl enable merge
systemctl start merge