# services/DEBIAN/postinst
#!/bin/sh

set -e
echo "Starting postinstall"

adduser --system --home /var/lib/merge merge || true

systemctl daemon-reload
systemctl enable merge.service
systemctl start merge.service

exit 0