# services/DEBIAN/prerm
#!/bin/sh

set -e
echo "Stopping the sort"

systemctl merge stop
systemctl disable merge
systemctl daemon-reload

echo "Merge Sort Service Stopped"