# services/DEBIAN/prerm
#!/bin/sh

set -e
echo "Stopping the sort"

systemctl stop merge.service
systemctl disable merge
systemctl daemon-reload

echo "Merge Sort Service Stopped"