# services/DEBIAN/prerm
#!/bin/sh

set -e
echo "Stopping the sort"

systemctl stop merge.service || true
systemctl disable merge.service || true
systemctl daemon-reload

echo "Merge Sort Service Stopped"